# docker-devserver/bootscript.service
[Unit]
After=mysqld.service

[Service]
Type=oneshot
ExecStart=/opt/bootscript.sh

[Install]
WantedBy=default.target
